# scripts/tempmon.sh
#!/bin/sh /etc/rc.common
#Openwrt on Raspberry Pi 4
while [ true ]
do
    # Code to run
        clear
        pi-cpu
        sleep 1
done